# airootfs/home/liveuser/.config/autostart/spicetheme.desktop
[Desktop Entry]
Name=Spice Theme
Exec=/usr/bin/spicetheme
Icon=utilities-terminal
Type=Application
Categories=System;Utility;
